# files/mysql_forwarder.service
[Unit]
Description=Forward MySQL to IPv4
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/socat TCP4-LISTEN:3306,fork TCP6:[<%= @ipaddress6 %>]:3306
ExecStop=/bin/pkill -SIGTERM -f /usr/bin/socat
KillMode=control-group
Restart=on-failure
RestartPreventExitStatus=255
PIDFile=/run/mysql_forwarder.pid
RestartSec=10s

[Install]
WantedBy=multi-user.target
